# scripts/crontab
# S3 Uploader cron schedule
# Install: crontab scripts/crontab
# Edit PROJ_DIR below to match your installation path.

# Interpret all schedules below in UTC (server may be in a non-UTC tz).
CRON_TZ=UTC
PROJ_DIR=/home/kartik/cardano-dune-analytics

# Daily pipeline with retry polling: starts at 02:00 UTC, retries every 10min until new files found
# Runs: dune (13 exporters) → drep_profile → off_chain_pool_data → drep_dist_enriched
# Sends Telegram notification on completion or failure
0 2 * * * $PROJ_DIR/scripts/daily-pipeline.sh >> $PROJ_DIR/logs/cron-daily.log 2>&1

# Smart contract registry: daily at 06:00 UTC (independent, checks GitHub for updates)
0 6 * * * $PROJ_DIR/scripts/run-exporter.sh --external smart_contract_registry >> $PROJ_DIR/logs/cron-registry.log 2>&1

# Asset data: every 2 hours (independent, fetches from Minswap API)
0 */2 * * * $PROJ_DIR/scripts/run-exporter.sh --external asset_data >> $PROJ_DIR/logs/cron-asset.log 2>&1
